# rabbitmq_hass/rootfs/etc/s6-overlay/s6-rc.d/rabbitmq/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# start rabbitmq-server
# ==============================================================================

set -e

bashio::log.info "Starting rabbitmq..."
exec s6-setuidgid rabbitmq /usr/sbin/rabbitmq-server
bashio::log.info "rabbitmq started!"
