# clickhouse_hass/rootfs/etc/s6-overlay/s6-rc.d/clickhouse-server/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Init postfix - generate config from options
# ==============================================================================

set -e

bashio::log.info "Starting clickhouse-server..."
exec /usr/bin/clickhouse-server
bashio::log.info "clickhouse-server started!"